# launchd/com.thoth.vaultkeeper-heartbeat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.thoth.vaultkeeper-heartbeat</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/risingtidesdev/.nvm/versions/node/v22.17.0/bin/node</string>
        <string>/Users/risingtidesdev/dev/Thoth/scripts/vaultkeeper-heartbeat.js</string>
        <string>--quiet</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/risingtidesdev/.nvm/versions/node/v22.17.0/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/risingtidesdev</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/risingtidesdev/dev/Thoth</string>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>StandardOutPath</key>
    <string>/Users/risingtidesdev/dev/Thoth/logs/vaultkeeper-heartbeat.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/risingtidesdev/dev/Thoth/logs/vaultkeeper-heartbeat-err.log</string>
</dict>
</plist>
