# launchd/LAUNCHD_PREFIX.living-dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>{{LAUNCHD_PREFIX}}.living-dashboard</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{NODE_PATH}}</string>
        <string>{{VAULT_PATH}}/scripts/sync-living-dashboard.js</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>{{NODE_PATH}}/../:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>{{VAULT_PATH}}</string>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>StandardOutPath</key>
    <string>{{VAULT_PATH}}/logs/living-dashboard-out.log</string>

    <key>StandardErrorPath</key>
    <string>{{VAULT_PATH}}/logs/living-dashboard-err.log</string>
</dict>
</plist>
